# deploy/systemd/centos9/multisoft-next.service
[Unit]
Description=Multisoft Next.js (CentOS 9)
After=network.target

[Service]
Type=simple
User=multisoft
Group=multisoft
WorkingDirectory=/opt/multisoft/multisoft_informes/frontend-next
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=AUTH_DB_HOST=localhost
Environment=AUTH_DB_PORT=5432
Environment=AUTH_DB_NAME=multisoft_informes
Environment=AUTH_DB_USER=multisoft_user
Environment=AUTH_DB_PASSWORD=Socio123!
Environment=NEXT_PUBLIC_NODE_API_URL=http://127.0.0.1:3000/api/v1
ExecStart=/bin/bash -lc 'source /home/multisoft/.nvm/nvm.sh && nvm use 20.20.2 >/dev/null && npm run start'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
